# recipes-connectivity/hostapd/hostapd/accesspoint.rules
SUBSYSTEM=="ieee80211", KERNEL=="phy*", ACTION=="add", RUN+="/usr/sbin/iw phy phy%n interface add uap%n type __ap"
